# pi/ansible/roles/panic-button/files/chimebox-panic-daemon.service
[Unit]
Description=chimebox panic-button daemon (kernel-level keystroke watcher)
Documentation=man:chimebox-panic-daemon(8)
# We need /dev/input/event* populated. multi-user.target ensures
# udev has settled. We don't depend on X.
After=multi-user.target
After=systemd-udevd.service
# Crash-loop protection: if the daemon dies more than 10 times in
# 60 seconds, give up. (StartLimit* keys belong in [Unit], not
# [Service] -- systemd warns at parse time if misplaced.)
StartLimitIntervalSec=60
StartLimitBurst=10

[Service]
Type=simple
ExecStart=/usr/local/sbin/chimebox-panic-daemon
Restart=always
RestartSec=2

# Hardening. The daemon needs:
#   - Read access to /dev/input/event* (already satisfied by
#     running as root)
#   - Permission to fork/exec the action scripts
#   - For force-reset: send signals to BasiliskII (CAP_KILL)
#   - For kid-reset: read+write /home/chimebox/chimebox (snapshot
#     dir + System.dsk) and chown the restored disk image
#     (CAP_CHOWN). This is why ProtectHome is off and why
#     CAP_CHOWN is in the bounding set.
# It does NOT need: arbitrary filesystem writes, network,
# new privileges, etc.
NoNewPrivileges=true
ProtectSystem=strict
# ProtectHome is intentionally OFF: kid-reset needs to read the
# snapshot dir and overwrite System.dsk, both under /home.
ProtectHome=false
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictNamespaces=true
RestrictRealtime=true
LockPersonality=true
MemoryDenyWriteExecute=true
# CapabilityBoundingSet limits what the daemon (and its action
# scripts) can do.
#   CAP_KILL            -- chimebox-force-reset / chimebox-emergency-stop
#                          send signals to BasiliskII
#   CAP_CHOWN           -- chimebox-reset chowns the restored
#                          System.dsk back to chimebox:chimebox
#   CAP_DAC_READ_SEARCH -- traverse 0750 chimebox-owned dirs as root
#                          (the daemon runs as root but, with a
#                          restricted bounding set, can't bypass
#                          DAC checks on dirs it doesn't own/group)
#   CAP_DAC_OVERRIDE    -- read snapshots and write System.dsk
#                          owned by chimebox:chimebox
#   CAP_FOWNER          -- mv across ownership for the .tmp -> final
#                          rename in chimebox-reset
CapabilityBoundingSet=CAP_KILL CAP_CHOWN CAP_DAC_READ_SEARCH CAP_DAC_OVERRIDE CAP_FOWNER
AmbientCapabilities=CAP_KILL CAP_CHOWN CAP_DAC_READ_SEARCH CAP_DAC_OVERRIDE CAP_FOWNER

# Logging goes to stderr (journal) and syslog (logger -t).
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
